# travis/html5-validator-service.conf
description "Validator.nu HTML5 Validator"

start on runlevel [234]
stop on runlevel [0156]

chdir /usr/share/html5-validator
exec python build/build.py --control-port=8889 run
respawn